# src/build/deb/etc/systemd/system/dragoon.service
[Unit]
Description=Dragoon video transcoder service

[Service]
User=dragoon
Group=dragoon
ExecStart=/usr/bin/java -Dlog4j.configurationFile=/opt/dragoon/log4j2.xml -jar '/opt/dragoon/dragoon.jar' -c '/etc/dragoon/config.toml'
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
